# config/systemd/backup-automation.service
[Unit]
Description=OpenStack Backup Automation
Documentation=https://github.com/example/openstack-backup-automation
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=backup
Group=backup
WorkingDirectory=/var/lib/backup-automation
Environment=CONFIG_FILE=/etc/backup-automation/config.yaml
ExecStart=/usr/local/bin/openstack-backup-automation run
StandardOutput=journal
StandardError=journal

# Security settings
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/backup-automation
ReadOnlyPaths=/etc/backup-automation
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true

[Install]
WantedBy=multi-user.target