# config/systemd/backup-automation.timer
[Unit]
Description=Run OpenStack Backup Automation every 15 minutes
Documentation=https://github.com/example/openstack-backup-automation
Requires=backup-automation.service

[Timer]
OnCalendar=*:0/15
Persistent=true
RandomizedDelaySec=60
AccuracySec=1min

[Install]
WantedBy=timers.target